# meta-xt-control-domain/recipes-guest/doma/files/doma-unpause.service
[Unit]
Description=Android VM unpause service

Requires=backend-ready@domd-qemu-monitor.service
After=backend-ready@domd-qemu-monitor.service

[Service]
Type=oneshot

ExecStart=/usr/sbin/xl unpause DomA

Restart=on-failure
RestartSec=5
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
